# supervisord-socat.conf
[program:socat]
command=/start-socat.sh
numprocs=1
autostart=true
autorestart=true